# rw_fs_stub/nomnt/etc/inittab
::sysinit:/etc/rc.d/start
ttyS0::respawn:-/bin/sh
